# systemd/sysd/ax25dev.path
[Path]
PathExists=/tmp/kisstnc

[Install]
WantedBy=multi-user.target
